# rc.d/maacom.in
#!/bin/sh
#
# $Id$
#
# PROVIDE: maacom
# REQUIRE: LOGIN

. /etc/rc.subr

name="maacom"
rcvar="maacom_enable"

pidfile="@APP_RUNDIR@"/maacom.pid
command="@prefix@/sbin/${name}"
command_args=""
procname="@prefix@/sbin/${name}"
command_interpreter="@perl@"

load_rc_config ${name}

: ${maacom_enable:="NO"}

run_rc_command "$1"
#EOF
