# maacom.service.in
[Unit]
Description=maacom

[Service]
Type=forking
PIDFile=@APP_RUNDIR@/maacom.pid
ExecStart=@prefix@/sbin/maacom
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target

